# manager-config/manager_ossec.conf
<integration>
  <name>virustotal</name>
  <api_key>YOUR_VIRUSTOTAL_API_KEY</api_key> <group>sysmon</group>
  <alert_format>json</alert_format>
</integration>

<integration>
  <name>slack</name>
  <hook_url>YOUR_SLACK_WEBHOOK_URL</hook_url> <level>10</level>
  <alert_format>json</alert_format>
</integration>

<command>
  <name>win-ip-block</name>
  <executable>netsh.exe</executable>
  <expect>srcip</expect>
  <timeout_allowed>yes</timeout_allowed>
</command>

<active-response>
  <command>win-ip-block</command>
  <location>local</location>
  <rules_id>60122,60204</rules_id> <timeout>60</timeout>
</active-response>

<group name="local,syslog,sshd,">

  <rule id="100002" level="12">
    <if_sid>61603</if_sid> <field name="win.eventdata.originalFileName" type="pcre2">(?i)vssadmin\.exe</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)delete\s+shadows</field>
    <description>Ransomware Activity: Shadow Copies deletion attempt detected</description>
    <mitre>
      <id>T1490</id>
    </mitre>
  </rule>

  <rule id="100003" level="10">
    <if_sid>61613</if_sid> <field name="win.eventdata.targetObject" type="pcre2">(?i)CurrentVersion\\\\Run</field>
    <description>Persistence: Suspicious modification of Registry Run Key</description>
    <mitre>
      <id>T1547</id>
    </mitre>
  </rule>

</group>